# etc/init.d/connect-spkr
#!/sbin/openrc-run

# Part of tiny-radioCD project, originally designed & hand-crafted by macmpi
# check License at project home https://github.com/macmpi/tiny-radioCD


description="Connects Bluetooth speaker for MPD"

depend() {
	use bluealsa
	after mpd
	before tiny-radiocd
}

start()
{
	ebegin "Connecting Bluetooth speaker"
		/usr/local/bin/connect-spkr
	eend $?
}

