# etc/init.d/easy-setup
#!/sbin/openrc-run

# Part of tiny-radioCD project, originally designed & hand-crafted by macmpi
# check License at project home https://github.com/macmpi/tiny-radioCD


description="Setup BT and Wifi settings for tiny-radioCD"

depend() {
	use bluealsa
	before sshd
	before dropbear
	before inotifyd
	before mpd
	before acpid
	before connect-spkr
}

start()
{
	ebegin "Starting tiny-radioCD easy-setup"
		/usr/local/bin/easy-setup
	eend $?
}

stop()
{
	ebegin
	eend $?
}

