# etc/init.d/monitor-spkr
#!/sbin/openrc-run

# Part of tiny-radioCD project, originally designed & hand-crafted by macmpi
# check License at project home https://github.com/macmpi/tiny-radioCD


description="Manages incoming bluetooth speaker connections"

name="Bluetooth speaker monitoring"
command="/usr/local/bin/monitor-spkr"
pidfile="/run/$RC_SVCNAME.pid"
command_background="yes"

depend() {
	use logger
	need bluealsa
	use mpd
}

