# etc/init.d/tiny-radiocd
#!/sbin/openrc-run

# Part of tiny-radioCD project, originally designed & hand-crafted by macmpi
# check License at project home https://github.com/macmpi/tiny-radioCD


description="tiny-radioCD init"

depend() {
	use mpd
}

start()
{
	ebegin "Starting tiny-radioCD"
		/usr/local/bin/radiocd-init
	eend $?
}

stop()
{
	ebegin
	eend $?
}

